# rpi/services/sentryspike-inference.service
[Unit]
Description=SentrySpike Inference Service
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/sentryspike
ExecStart=/opt/sentryspike/venv/bin/python SentrySpike_Inference/inference_service.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
